# recipes-cockpit/cockpit-streambox-config/files/yocto/cockpit-streambox-config.socket
[Unit]
Description=StreamBox Config API Socket for Cockpit

[Socket]
ListenStream=9091

[Install]
WantedBy=sockets.target
